# packer/files/etc/cron.d/sync_ssh_keys_cron
*/5 * * * * root /bin/bash /usr/local/bin/sync_ssh_keys.sh
